# support/systemd/boarai-firmware.service
[Unit]
Description=boarAI Firmware

[Service]
ExecStart=/home/hsr/code/firmware/run-firmware.sh
Type=simple

[Install]
WantedBy=default.target
